# software/meta-package/debian/shepherd.shepherd-rpc.service
[Unit]
Description=shepherd zero rpc server

[Service]
ExecStart=/usr/local/bin/shepherd-sheep -vv rpc
Restart=on-failure
