# software/gps/chrony.path

[Unit]
Description= Triggers chrony service when PPS device appears

[Path]
PathExists=/dev/pps0

[Install]
WantedBy=multi-user.target