# software/ptp/debian/postinst
systemctl enable phc2sys
systemctl disable systemd-timesyncd
