# deploy/roles/sheep/files/shepherd-launcher.service
[Unit]
Description=shepherd launcher

[Service]
ExecStart=/usr/local/bin/shepherd-sheep -vv launcher
Restart=on-failure
