# deploy/roles/sheep/files/shepherd-gui.service
[Unit]
Description=Shepherd debug gui

[Service]
ExecStart=/usr/local/bin/shepherd-sheep gui
Restart=always
