# software/deploy/roles/ptp-slave/templates/ptp4l.path
[Unit]
Description=ptp4l path activation

[Path]
PathExists=/dev/ptp0

[Install]
WantedBy=multi-user.target
